# deploy/com.neotrix.proxy-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.neotrix.proxy-daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/3tfire/Downloads/neotrix/target/release/neotrix-proxy-daemon</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/neotrix-proxy-daemon.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/neotrix-proxy-daemon.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>RUST_LOG</key>
        <string>info</string>
    </dict>
</dict>
</plist>
